# rootfs/etc/services.d/dns/run
#!/usr/bin/with-contenv sh
exec go-dnsmasq --syslog --default-resolver --append-search-domains --hostsfile=/etc/hosts ${DNSMASQ_ARGS} &>/dev/null
