# rootfs/etc/services.d/syslog/run
#!/usr/bin/env sh
exec syslogd -n -D -O /dev/stdout